# file-service/file-service.service
[Unit]
Description=Yessal File Service (Recette)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/proprojects/yessal/file-service
ExecStart=/usr/bin/node src/server.js
Restart=on-failure
RestartSec=10

# Environment
Environment=NODE_ENV=production
Environment=PORT=4540

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=yessal-file-service

[Install]
WantedBy=multi-user.target